# banometer.service
[Unit]
Description=Banometer app service

[Service]
ExecStart=./scripts/start.sh